# web_rl.socket
# /etc/systemd/system/web_rl.socket

[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/web_rl.sock

[Install]
WantedBy=sockets.target
